# config/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=error

logfile_maxbytes=50MB
logfile_backups=10

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[eventlistener:logfilter]
command=/bin/bash -c "grep -v 'reaped unknown pid' || true"
events=PROCESS_LOG
buffer_size=1024

[program:cups]

command=/bin/bash -c "sleep 2 && exec /usr/sbin/cupsd -f"
autostart=true
autorestart=true
priority=10
startsecs=5
startretries=10
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/cups.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/cups_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

[program:printer_client]
command=/usr/bin/php /opt/websocket_printer/printer_client.php
directory=/opt/websocket_printer
autostart=true
autorestart=true
priority=20
startsecs=5
startretries=999
stopwaitsecs=10
environment=HOME="/tmp",USER="root",SHELL="/bin/bash"
stdout_logfile=/var/log/supervisor/printer_client.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/printer_client_err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
priority=5
startsecs=3
stdout_logfile=/var/log/supervisor/sshd.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/sshd_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3